# .local/.cpyroot/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --login-options skeignark --noclear %I $TERM
Type=idle
